# system/launchd/com.user.rufino-ingest-applehealth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.rufino-ingest-applehealth</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.claude/scripts/rufino-ingest-applehealth.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/logs/rufino/rufino-ingest-applehealth.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/logs/rufino/rufino-ingest-applehealth.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>RUFINO_VAULT_PATH</key>
        <string></string>
        <key>RUFINO_DISPLAY_NAME</key>
        <string></string>
    </dict>
</dict>
</plist>
